# snapserver-addon/rootfs/etc/services.d/avahi/run
#!/command/with-contenv bashio
# ==============================================================================
# Runs avahi-daemon
# ==============================================================================
bashio::log.info "Starting avahi-daemon..."

# Run avahi-daemon
exec avahi-daemon --no-drop-root --debug
